# chronia.service
[Unit]
Description=Chronia Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/Chronia

[Install]
WantedBy=multi-user.target
